# nginx-loadbalancer/www/etc/supervisord.d/update_html.ini
[program:update_html]
command=/usr/local/bin/update_html.sh
redirect_stderr=true
exitcodes=0
startsecs=0
stdout_logfile=/var/log/supervisor/%(program_name)s.log
