# nginx-loadbalancer/lb/etc/supervisord.d/lb_update.ini
[program:lb-update]
command=consul-template -consul localhost:8500 -template "/etc/consul-templates/lb.conf.ctmpl:/etc/nginx/conf.d/lb.conf:nginx -s reload"
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
